# mimo.service
[Unit]
Description=Mimo AI Flask Application
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/root/uncensored_llm
Environment=LLAMA_N_CTX=2048
Environment=LLAMA_N_THREADS=4
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/root/uncensored_llm
ExecStart=/usr/bin/python3 /root/uncensored_llm/start_production.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mimo-ai

# Memory optimization
MemoryLimit=4G
MemoryMax=6G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
